# monitor_temp.service
[Unit]
Description=Monitors core temperature and runs GPIO fan

[Service]
Type=simple
ExecStart=/opt/fan_control 55 5

